# .devcontainer/config/mysql.conf
[program:mysql]
command=/usr/bin/mysqld_safe --user=mysql
priority=1000
autostart=true
autorestart=true
user=root